# random/units/minecraftsv
#!/sbin/openrc-run
# Copyright 2024 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

#supervisor=supervise-daemon
name="Mein Kampf Server"
description=""
command_user="servergas"
command="/srv/servergas/.local/bin/startsv.sh"
procname="servergas-tmux-wrapper"
stop() {
	su servergas -l -c '/srv/servergas/.local/bin/bellanno.sh'
}

depend() {
	use dns
}
