# gentoo/amd_a10_7860k/etc/init.d/ngrok
#!/sbin/openrc-run
# Copyright 2024 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2
supervisor=supervise-daemon
USER=ngrok
name="ngrok daemon"
description=""
command_user="${USER}"
command="/srv/ngrok/bin/ngrok"
command_args="${ngrok_args} start --all --config /srv/ngrok/config/ngrok1.yml"

depend() {
	use dns
}
