# random/oldetc/runit/playit/run
#!/bin/sh -e
exec 2>&1
exec chpst -u playit /srv/playit/bin/playit --secret_path /srv/playit/playit.toml
